# buildroot-external/rootfs-overlay/etc/systemd/system/mycroft-messagebus.service.d/splash.conf
[Service]
ExecStartPre=-/bin/psplash-write "MSG Starting Mycroft A.I. messagebus service"
ExecStartPre=-/bin/sleep 2
