# buildroot-external/rootfs-overlay/base/usr/lib/systemd/system/rpcbind.socket
[Unit]
Description=RPCbind Server Activation Socket

[Socket]
ListenStream=/run/rpcbind.sock

[Install]
WantedBy=sockets.target
