# buildroot-external/rootfs-overlay/etc/systemd/system/spotifyd.service.d/splash.conf
[Service]
ExecStartPre=-/bin/psplash-write "MSG Starting Spotify connect"
