# buildroot-external/rootfs-overlay/etc/NetworkManager/dispatcher.d/wifi-toggle
#!/bin/sh

if [ "$1" = "eth0" ]; then
    case "$2" in
        up)
            nmcli radio wifi off
            ;;
        down)
            nmcli radio wifi on
            ;;
    esac
fi
