# buildroot-external/rootfs-overlay/usr/lib/systemd/system/connectivity-check.service
[Unit]
DefaultDependencies=no
After=nss-lookup.target
Before=network-online.target
Type=oneshot
RemainAfterExit=yes

[Service]
ExecStart=sh -c 'while ! ping -c 1 docker.io; do sleep 1; done'

[Install]
WantedBy=network-online.target
