# buildroot-external/rootfs-overlay/usr/lib/systemd/system/mycroft-splash-40.service
[Unit]
Description=Update Mycroft splash boot screen to 40%
Wants=NetworkManager.service
After=NetworkManager.service
DefaultDependencies=no

[Service]
ExecStart=-/usr/bin/psplash-write "PROGRESS 40" 

[Install]
WantedBy=multi-user.target
