# buildroot-external/rootfs-overlay/etc/systemd/system/snapserver.service.d/splash.conf
[Service]
ExecStartPre=-/bin/psplash-write "MSG Starting Snapcast server"
