# buildroot-external/rootfs-overlay/usr/lib/systemd/system/seeed-voicecard.service.d/splash.conf
[Service]
ExecStartPre=-/bin/psplash-write "MSG Preparing respeaker configuration"
