# buildroot-external/rootfs-overlay/usr/lib/systemd/system/prepare_system.service.d/splash.conf
[Service]
ExecStartPre=-/bin/psplash-write "MSG Preparing system"
