# buildroot-external/rootfs-overlay/etc/udev/rules.d/70-ap-interface.rules
SUBSYSTEM=="net", KERNEL=="wlan*", ACTION=="add", RUN+="/sbin/iw dev %k interface add ap%n type __ap"
