# buildroot-external/rootfs-overlay/usr/lib/systemd/system/mycroft.service.d/splash.conf
[Service]
ExecStartPre=-/bin/psplash-write "MSG Starting Mycroft A.I. software stack"
ExecStartPre=-/bin/sleep 2
