# buildroot-external/rootfs-overlay/usr/lib/systemd/user/mycroft.service
[Unit]
Description=Mycroft A.I. Software stack.

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=default.target
