# buildroot-external/rootfs-overlay/etc/systemd/system/firstboot.service.d/splash.conf
[Service]
ExecStartPre=-/bin/psplash-write "MSG Preparing for first boot"
