# buildroot-external/rootfs-overlay/usr/lib/systemd/system/mycroft-voice.service.d/splash.conf
[Service]
ExecStartPre=-/bin/psplash-write "MSG Starting Mycroft A.I. voice service"
ExecStartPre=-/bin/sleep 2
