# buildroot-external/rootfs-overlay/usr/lib/systemd/system/mycroft-messagebus.service.d/splash.conf
[Service]
ExecStartPre=-/bin/psplash-write "PROGRESS 90"
