# buildroot-external/rootfs-overlay/usr/lib/systemd/system/hostname.service.d/splash.conf
[Service]
ExecStartPre=-/bin/psplash-write "MSG Setting up hostname"
