# buildroot-external/rootfs-overlay/etc/systemd/system/mycroft-audio.service.d/splash.conf
[Service]
ExecStartPre=-/bin/psplash-write "MSG Starting Mycroft A.I. audio service"
ExecStartPre=-/bin/sleep 2
