# buildroot-external/rootfs-overlay/usr/lib/systemd/system/mycroft-enclosure.service.d/splash.conf
[Service]
ExecStartPre=-/bin/psplash-write "MSG Starting Mycroft A.I. enclosure service"
ExecStartPre=-/bin/sleep 2
