# buildroot-external/rootfs-overlay/usr/lib/systemd/system/mycroft-enclosure.service.d/splash.conf
[Service]
ExecStartPre=-/usr/bin/sudo /bin/psplash-write "PROGRESS 95"
