# buildroot-external/rootfs-overlay/etc/systemd/system/mycroft-skills.service.d/splash.conf
[Service]
ExecStartPre=-/bin/psplash-write "MSG Starting Mycroft A.I. skills service"
ExecStartPre=-/bin/sleep 2
