# buildroot-external/package/respeaker/seeed-voicecard.service
[Unit]
Description=Seeed Voicecard service
After=alsa-restore.service
 
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/bin/psplash-write "MSG Preparing respeaker configuration"
ExecStart=/usr/bin/seeed-voicecard
User=root
 
[Install]
WantedBy=sysinit.target
