# buildroot-external/board/ovos/raspberrypi/rootfs-overlay/base/etc/udev/rules.d/98-backlight-permissions.rules
SUBSYSTEM=="backlight",RUN+="/bin/chmod 666 /sys/class/backlight/%k/brightness /sys/class/backlight/%k/bl_power"
