# buildroot-external/rootfs-overlay/home/ovos/.config/systemd/user/pipewire.service.d/openvoiceos.conf
[Service]
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=49
