# buildroot-external/rootfs-overlay/etc/systemd/system/snapclient.service.d/splash.conf
[Service]
ExecStartPre=-/bin/psplash-write "MSG Starting Snapcast client"
