# buildroot-external/package/mycroft-service/mycroft.service
[Unit]
Description=Mycroft A.I. Software stack.
After=pulseaudio.service NetworkManager.service systemd-timesyncd.service

[Service]
User=mycroft
WorkingDirectory=/home/mycroft
Type=oneshot
ExecStartPre=-/usr/bin/sj201-reset-led
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
