# buildroot-external/package/mycroft-service/mycroft.service
[Unit]
Description=Mycroft A.I. Software stack.
After=pulseaudio.service

[Service]
User=mycroft
WorkingDirectory=/home/mycroft
Type=oneshot
ExecStartPre=/bin/psplash-write "MSG Starting Mycroft A.I. software stack"
ExecStartPre=/bin/sleep 2
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
