# buildroot-external/rootfs-overlay/usr/lib/systemd/system/mycroft-splash-start.service.d/splash.conf
[Service]
ExecStartPost=-/bin/sleep 2
ExecStartPost=-/bin/psplash-write "PROGRESS 10"
