# buildroot-external/rootfs-overlay/usr/lib/systemd/system/getty@tty1.service.d/override.conf
[Service]
Type=simple
ExecStart=
ExecStart=-/usr/sbin/agetty --autologin mycroft --noclear %I $TERM
